# scrutiny.cron
*/15 * * * * /home/$CONTAINERUSER/scrutiny/collector/scrutiny-collector-metrics-linux-amd64 run --api-endpoint "http://0.0.0.0:8080"
