# scrutiny.cron
*/15 * * * * root /scrutiny/collector/scrutiny-collector-metrics-linux-${ARCH} run --api-endpoint ${SCRUTINY_API_URL}
