# scrutiny.cron
*/15 * * * * . /etc/profile; /home/$CONTAINERUSER/scrutiny/collector/scrutiny-collector-metrics-linux-amd64 run --api-endpoint "http://localhost:8080"
