# scrutiny.cron
*/1 * * * * root  /scrutiny/collector/collector.sh
