# scrutiny.cron
*/15 * * * * /scrutiny/collector/scrutiny-collector-metrics-linux-${ARCH} run --api-endpoint "http://app:8080"
